# deploy/systemd/aisurv-consumer.service
[Unit]
Description=AI Surveillance Consumer Workers
After=network.target redis.service postgresql.service

[Service]
User=aisurv
WorkingDirectory=/opt/aisurv/backend
ExecStart=/opt/aisurv/venv/bin/python -m app.consumers.runner
Restart=always
RestartSec=3
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target

